# crontab
*/20 * * * * /home/rpi/photobooth-thing/wakeup-printer.sh
